# config/etc/crontab
# Speed up UEFI boot
@reboot root:wheel /usr/sbin/efibootmgr -n -b x